# etc/cron.d/defroute_chk
*/1 * * * * root /etc/network/scripts/defroute_chk
